# turbo/apps/runner/scripts/deploy/vsock-agent.service
[Unit]
Description=VM0 Vsock Agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/vm0-agent/vsock-agent.py
Restart=always
RestartSec=1
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
